# timelimit-ui/rootfs/etc/s6-overlay/s6-rc.d/frontend-ui/run
#!/usr/bin/with-contenv bashio
cd /usr/bin
bashio::log.info "Starting Frontend Webserver..."
# De -u vlag zorgt dat logs DIRECT verschijnen
exec python3 -u frontend.py