# timelimit-ui/rootfs/etc/s6-overlay/s6-rc.d/frontend-ui/run
#!/usr/bin/with-contenv bashio

# Ga naar de map waar alle bestanden staan
cd /usr/bin

bashio::log.info "Starting Frontend Webserver..."

# Start de server vanuit de huidige map
exec python3 frontend.py