# rootfs/etc/services.d/timelimit-ui/finish
#!/usr/bin/execlineb -S0
foreground { echo "[timelimit-ui] Backend service stopped" }
exit 0