# timelimit-ui/rootfs/etc/services.d/timelimit-ui/run
#!/usr/bin/with-contenv bash
echo "TimeLimit UI add-on started"
sleep infinity