# timelimit-ui/rootfs/etc/services.d/test/run
#!/usr/bin/with-contenv bash
echo "Service started"
sleep 999999