# rootfs/etc/services.d/backend/finish
#!/usr/bin/with-contenv bashio
echo "[timelimit-ui] Backend service stopped"