# timelimit-ui/rootfs/etc/s6-overlay/s6-rc.d/backend/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Timelimit backend communication check..."

# Voer het python script uit
python3 /usr/bin/check_server.py