# rootfs/etc/s6-overlay/s6-rc.d/example-service/run
#!/usr/bin/with-contenv bashio
echo "Start de Timelimit UI processen..."
# Hier komt later je echte commando, bijv: python3 -m timelimit
sleep infinity