# timelimit-ui/rootfs/etc/s6-overlay/s6-rc.d/web-server/run
#!/usr/bin/with-contenv bashio
cd /usr/bin
bashio::log.info "Starting TimeLimit UI Web Server..."
# De -u vlag zorgt dat logs DIRECT verschijnen
exec python3 -u web_server.py